# etc/cron.hourly/auto-whitelist
#!/bin/bash

/usr/local/bin/auto-whitelist >/dev/null
