# etc/cron.hourly/auto-welcomelist
#!/bin/bash

/usr/local/bin/auto-welcomelist >/dev/null
