# etc/cron.daily/dovecot-archive
#!/bin/bash

/usr/local/bin/dovecot-archive >/dev/null
