# etc/cron.weekly/rotate
#!/bin/bash

/usr/local/bin/rotate >/dev/null
